# files/etc/sv/fcgiwrap/run
#!/bin/sh

[ -e "/run/fcgiw.sock" ] && rm /run/fcgiw.sock

fcgiwrap -s unix:/run/fcgiw.sock
