# files/etc/sv/php-fpm8.2/run
#!/bin/sh
[ -r conf ] && . ./conf
[ -d /run/php ] || mkdir /run/php
exec php-fpm8.2 --nodaemonize ${OPTS}
